# supervisor/homaton.conf
[program:homaton]
command=/usr/bin/python homaton.py
user=root
stderr_logfile=/var/log/supervisor/homaton-err.log
stdout_logfile=/var/log/supervisor/homaton-out.log
directory=/home/pi/homaton/web/
